# autostart/.startvicinae.sh.desktop
[Desktop Entry]
Exec=/home/arshan/.startvicinae.sh
Icon=application-x-shellscript
Name=.startvicinae.sh
Type=Application
X-KDE-AutostartScript=true
